# crontab
#每分钟备份一次
*/1 * * * * /root/xtrbackup/xtrbackup_backup.sh > /data/backups/backup-`date "+\%G:\%m:\%d-\%H:\%M:\%S"`.log 2>&1
